# deploy/systemd/pipeworks-policy-workbench.service
[Unit]
Description=PipeWorks Policy Workbench
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pipeworks
Group=pipeworks
WorkingDirectory=/srv/work/pipeworks/repos/pipeworks-policy-workbench
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/pipeworks/policy-workbench/policy-workbench.env
ExecStart=/srv/work/pipeworks/venvs/pw-policy-workbench/bin/pw-policy serve --host ${PW_POLICY_SERVICE_HOST} --port ${PW_POLICY_SERVICE_PORT}
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
